# dollyburns/dollyburns.cron.sample
# /etc/cron.d/dollyburns: crontab entry for the dollyburns package

RUBY=/usr/bin/ruby
DOLLYBURNS=/usr/bin/dollyburns.rb

00 07   * * *   root    $RUBY $DOLLYBURNS
